# crontab
0 0 * * * ec2-user /Users/cperry/Developer/GitRepos/tools/aws-missing-tools/ec2-automate-backup/ec2-automate-backup-awscli.sh -s tag -t 'Backup,Values=true' -k 7 -p -n
